# src/systemd/user/wireplumber@.service.in
[Unit]
Description=Multimedia Service Session Manager (%i)
After=pipewire.service
BindsTo=pipewire.service
AssertPathExists=@WP_CONFIG_DIR@/%i.conf
Conflicts=pipewire-media-session.service

[Service]
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
Type=simple
ExecStart=@WP_BINARY@ -c %i.conf
Restart=on-failure
Slice=session.slice

[Install]
WantedBy=pipewire.service
